# puppet-manifests/debian/trixie/src/modules/platform/templates/ptp4l-instance.service.erb
[Unit]
Description=Precision Time Protocol (PTP) service
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=all
EnvironmentFile=-<%= @ptp_options_dir %>/ptpinstance/ptp4l-instance-%i
ExecStart=/usr/local/bin/ptp-instance-notify.sh ptp4l %i <%= @ptp_conf_dir %>
ExecStopPost=/bin/rm -f /var/run/ptp4l-%i.pid
Restart=on-failure
RestartPreventExitStatus=SIGTERM SIGINT SIGKILL SIGABRT 255
TimeoutStartSec=10
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=70

[Install]
WantedBy=multi-user.target
WantedBy=ptp4l-instances.target
