# puppet-manifests/debian/trixie/src/modules/platform/templates/ptp4l-instances.target.erb
[Unit]
Description=All ptp4l instances started
After=network-online.target

[Install]
WantedBy=multi-user.target
