# puppet-manifests/debian/trixie/src/modules/platform/templates/synce4l-instances.target.erb
[Unit]
Description=All synce4l instances started
After=ptp4l-instances.target

[Install]
WantedBy=multi-user.target
