# puppet-manifests/debian/trixie/src/modules/platform/templates/ts2phc-instances.target.erb
[Unit]
Description=All ts2phc instances started
After=ptp4l-instances.target

[Install]
WantedBy=multi-user.target
