# puppet-manifests/src/modules/platform/templates/gpsd.service.erb
[Unit]
Description=GPS (Global Positioning System) Daemon
Requires=gpsd.socket
# Needed with chrony SOCK refclock
After=chronyd.service

[Service]
Type=forking
EnvironmentFile=-<%= @ptp_options_dir %>/ptpinstance/monitoring-ptp
ExecStart=/usr/sbin/gpsd $GPSD_OPTIONS $DEVICES
ExecStartPost=/bin/bash -c 'echo $MAINPID > /var/run/gpsd.pid'
ExecStopPost=/bin/rm -f /var/run/gpsd.pid
Restart=on-failure
RestartPreventExitStatus=SIGTERM SIGINT SIGKILL SIGABRT 255

[Install]
WantedBy=multi-user.target
Also=gpsd.socket
