# puppet-manifests/src/modules/platform/templates/k8splatform.slice.erb
[Unit]
Description=Kubernetes Platform slice
Before=slices.target

[Slice]
MemoryAccounting=true
CPUAccounting=true
CPUShares=<%= @k8splatform_cpushares %>
