# puppet-manifests/src/modules/platform/templates/phc2sys-instance.service.erb
[Unit]
Description=Synchronize system clock or PTP hardware clock (PHC)
After=ptp4l.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/phc2sys-instance
ExecStart=/usr/sbin/phc2sys -f /etc/phc2sys-%i.conf $OPTIONS
ExecStartPost=/bin/bash -c 'echo $MAINPID > /var/run/phc2sys-%i.pid'
ExecStopPost=/bin/rm -f /var/run/phc2sys-%i.pid
Restart=on-failure
RestartPreventExitStatus=SIGTERM SIGINT SIGKILL SIGABRT 255
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=70

[Install]
WantedBy=multi-user.target